# src/Domain/Resources/docker/db/supervisor.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart=true
